# cloud-automation-concepts/Assignment 1 - AWS Basics/config/cloudshirt.service
[Unit]
Description=CloudShirt ASP.NET Core Web App
After=network.target

[Service]
WorkingDirectory=/opt/cloudshirt/publish
ExecStart=/usr/bin/dotnet /opt/cloudshirt/publish/Web.dll
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=cloudshirt
User=ec2-user

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
Environment=DOTNET_ROOT=/usr/share/dotnet

# Verbindingsstrings en app-config worden ingelezen vanuit .env (REQ-04)
EnvironmentFile=/opt/cloudshirt/.env

[Install]
WantedBy=multi-user.target
